# novnc-conf/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc 6080 %(ENV_VNC_HOST)s:%(ENV_VNC_PORT)s
autorestart=true
